# systemd/kalseer-daily.timer
[Unit]
Description=Run the Kalseer morning brief daily at 08:00 (catches up after boot)

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
